# supervisord.conf
# =============================================================================
# AI STRATEGY HUB - SUPERVISORD CONFIGURATION
# =============================================================================
# Gestisce tutti i processi nel container monolitico:
# - PostgreSQL 15
# - Redis 7
# - FastAPI Backend
# - Next.js Frontend
# - Nginx (routing interno)
# - Background Tasks Worker
#
# Autore: Claude per Davide
# Data: 2026-01-15
# =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/app/logs
loglevel=info

# -----------------------------------------------------------------------------
# PostgreSQL 15
# -----------------------------------------------------------------------------
[program:postgresql]
command=/usr/lib/postgresql/15/bin/postgres -D /var/lib/postgresql/15/main
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/app/logs/postgresql.out.log
stderr_logfile=/app/logs/postgresql.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5

# Importante: PostgreSQL deve essere pronto prima degli altri servizi
startsecs=10
stopwaitsecs=30

# -----------------------------------------------------------------------------
# Redis 7
# -----------------------------------------------------------------------------
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf --daemonize no
user=redis
autostart=true
autorestart=true
priority=20
stdout_logfile=/app/logs/redis.out.log
stderr_logfile=/app/logs/redis.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5

startsecs=5
stopwaitsecs=10

# -----------------------------------------------------------------------------
# FastAPI Backend
# -----------------------------------------------------------------------------
[program:backend]
command=/usr/local/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
directory=/app/backend
user=root
autostart=true
autorestart=true
priority=30
environment=PYTHONPATH="/app/backend",PYTHONUNBUFFERED="1"
stdout_logfile=/app/logs/backend.out.log
stderr_logfile=/app/logs/backend.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10

# Aspetta che PostgreSQL e Redis siano pronti
startsecs=15
stopwaitsecs=30

# Restart automatico in caso di crash
startretries=3

# -----------------------------------------------------------------------------
# Next.js Frontend
# -----------------------------------------------------------------------------
[program:frontend]
command=/usr/bin/npm start
directory=/app/frontend
user=root
autostart=true
autorestart=true
priority=40
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0"
stdout_logfile=/app/logs/frontend.out.log
stderr_logfile=/app/logs/frontend.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10

startsecs=10
stopwaitsecs=20

# -----------------------------------------------------------------------------
# Nginx (Reverse Proxy Interno)
# -----------------------------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=50
stdout_logfile=/app/logs/nginx.out.log
stderr_logfile=/app/logs/nginx.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10

# Nginx deve partire dopo backend e frontend
startsecs=5
stopwaitsecs=10

# -----------------------------------------------------------------------------
# Background Tasks Worker
# -----------------------------------------------------------------------------
[program:background_worker]
command=/usr/local/bin/python -m app.workers.task_runner
directory=/app/backend
user=root
autostart=true
autorestart=true
priority=60
environment=PYTHONPATH="/app/backend",PYTHONUNBUFFERED="1"
stdout_logfile=/app/logs/worker.out.log
stderr_logfile=/app/logs/worker.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile_backups=10

# Worker può partire dopo il backend
startsecs=10
stopwaitsecs=30

# Restart automatico
startretries=3

# -----------------------------------------------------------------------------
# Group: Core Services
# -----------------------------------------------------------------------------
# Raggruppa i servizi core per gestirli insieme
[group:core]
programs=postgresql,redis
priority=10

# -----------------------------------------------------------------------------
# Group: Application Services
# -----------------------------------------------------------------------------
[group:application]
programs=backend,frontend,background_worker
priority=30

# -----------------------------------------------------------------------------
# Group: Web Server
# -----------------------------------------------------------------------------
[group:webserver]
programs=nginx
priority=50

# -----------------------------------------------------------------------------
# Event Listeners (opzionali - per monitoring avanzato)
# -----------------------------------------------------------------------------
# [eventlistener:crashmail]
# command=/path/to/crashmail
# events=PROCESS_STATE_EXITED

# -----------------------------------------------------------------------------
# Configurazione Log Rotation (gestito da logrotate esterno)
# -----------------------------------------------------------------------------
# I log sono limitati tramite maxbytes/backups sopra
# Per rotation più sofisticata, usare logrotate del sistema

# =============================================================================
# NOTES
# =============================================================================
# - Priority: numeri più bassi = avvio prima
# - PostgreSQL (10) e Redis (20) partono per primi
# - Backend (30) e Frontend (40) partono dopo
# - Nginx (50) parte per ultimo (quando tutto è pronto)
# - Worker (60) può partire in parallelo con gli altri
#
# Comandi utili:
# - supervisorctl status          : Stato di tutti i processi
# - supervisorctl restart backend : Riavvia solo il backend
# - supervisorctl stop all        : Ferma tutti i servizi
# - supervisorctl start all       : Avvia tutti i servizi
# =============================================================================
